# systemd/battery_warn.service
[Unit]
Description=Display a battery warning

[Service]
Type=oneshot
ExecStart=/usr/local/bin/battery_warn
